# onecloud-pluginipk/usr/lib/opkg/info/plugin-boxweb.list
/miner.plugin-boxweb.ipk/conf/zlog.conf
/miner.plugin-boxweb.ipk/conf/param
/miner.plugin-boxweb.ipk/makeipkversion
/miner.plugin-boxweb.ipk/usr/sbin/fcgiwrap
/miner.plugin-boxweb.ipk/usr/sbin/nginx
/miner.plugin-boxweb.ipk/usr/bin/cgi-fcgi
/miner.plugin-boxweb.ipk/start.sh
/miner.plugin-boxweb.ipk/etc/init.d/S91nginx
/miner.plugin-boxweb.ipk/etc/nginx/uwsgi_params
/miner.plugin-boxweb.ipk/etc/nginx/scgi_params
/miner.plugin-boxweb.ipk/etc/nginx/nginx.conf.default
/miner.plugin-boxweb.ipk/etc/nginx/koi-win
/miner.plugin-boxweb.ipk/etc/nginx/mime.types.default
/miner.plugin-boxweb.ipk/etc/nginx/scgi_params.default
/miner.plugin-boxweb.ipk/etc/nginx/uwsgi_params.default
/miner.plugin-boxweb.ipk/etc/nginx/fastcgi.conf.default
/miner.plugin-boxweb.ipk/etc/nginx/mime.types
/miner.plugin-boxweb.ipk/etc/nginx/koi-utf
/miner.plugin-boxweb.ipk/etc/nginx/nginx.conf
/miner.plugin-boxweb.ipk/etc/nginx/fastcgi_params
/miner.plugin-boxweb.ipk/etc/nginx/fastcgi_params.default
/miner.plugin-boxweb.ipk/etc/nginx/win-utf
/miner.plugin-boxweb.ipk/etc/nginx/fastcgi.conf
/miner.plugin-boxweb.ipk/etc/logrotate.d/nginx
/miner.plugin-boxweb.ipk/version
/miner.plugin-boxweb.ipk/python_scripts/device_info.py
/miner.plugin-boxweb.ipk/python_scripts/reboot.py
/miner.plugin-boxweb.ipk/python_scripts/set_network.py
/miner.plugin-boxweb.ipk/python_scripts/diag.sh
/miner.plugin-boxweb.ipk/python_scripts/get_network.py
/miner.plugin-boxweb.ipk/python_scripts/timestamp_set.py
/miner.plugin-boxweb.ipk/python_scripts/reset.py
/miner.plugin-boxweb.ipk/python_scripts/activate.py
/miner.plugin-boxweb.ipk/lib/libfcgi++.so.0.0.0
/miner.plugin-boxweb.ipk/lib/libfcgi.so.0.0.0
/miner.plugin-boxweb.ipk/nginx_html/diag.cgi
/miner.plugin-boxweb.ipk/nginx_html/set_network.cgi
/miner.plugin-boxweb.ipk/nginx_html/reset.cgi
/miner.plugin-boxweb.ipk/nginx_html/web/umi.67d12b15.js
/miner.plugin-boxweb.ipk/nginx_html/web/vendors.610e5b27.async.js
/miner.plugin-boxweb.ipk/nginx_html/web/5.5328437d.async.js
/miner.plugin-boxweb.ipk/nginx_html/web/icons/icon-128x128.png
/miner.plugin-boxweb.ipk/nginx_html/web/icons/icon-192x192.png
/miner.plugin-boxweb.ipk/nginx_html/web/icons/icon-512x512.png
/miner.plugin-boxweb.ipk/nginx_html/web/3.b65bc269.async.js
/miner.plugin-boxweb.ipk/nginx_html/web/6.c63e30b2.chunk.css
/miner.plugin-boxweb.ipk/nginx_html/web/7.ca738295.async.js
/miner.plugin-boxweb.ipk/nginx_html/web/6.72a55244.async.js
/miner.plugin-boxweb.ipk/nginx_html/web/favicon.png
/miner.plugin-boxweb.ipk/nginx_html/web/2.3dadf6c0.async.js
/miner.plugin-boxweb.ipk/nginx_html/web/vendors.bf7359fe.chunk.css
/miner.plugin-boxweb.ipk/nginx_html/web/4.e88ead53.chunk.css
/miner.plugin-boxweb.ipk/nginx_html/web/4.0ccfef46.async.js
/miner.plugin-boxweb.ipk/nginx_html/web/5.b752e68e.chunk.css
/miner.plugin-boxweb.ipk/nginx_html/web/3.933e82a3.chunk.css
/miner.plugin-boxweb.ipk/nginx_html/web/umi.2671cf43.css
/miner.plugin-boxweb.ipk/nginx_html/web/index.html
/miner.plugin-boxweb.ipk/nginx_html/web/2.73f2fe96.chunk.css
/miner.plugin-boxweb.ipk/nginx_html/get_network.cgi
/miner.plugin-boxweb.ipk/nginx_html/reboot.cgi
/miner.plugin-boxweb.ipk/nginx_html/favicon.png
/miner.plugin-boxweb.ipk/nginx_html/activate.cgi
/miner.plugin-boxweb.ipk/nginx_html/device_info.cgi
/miner.plugin-boxweb.ipk/nginx_html/index.html
/miner.plugin-boxweb.ipk/nginx_html/timestamp_set.cgi
/miner.plugin-boxweb.ipk/bin/boxweb_monitor
/miner.plugin-boxweb.ipk/bin/run.sh
/miner.plugin-boxweb.ipk/bin/boxweb
/miner.plugin-boxweb.ipk/bin/readme
/miner.plugin-boxweb.ipk/stop.sh
