# etc/crontabs/root
0 3 * * * cd / && run-parts /etc/cron.daily  # 每天凌晨3点执行

